# roles/managed_host/files/fstrim.timer
[Unit]
Description=regularly discard unused blocks on file systems

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
