# roles/managed_node/files/fwupd-refresh.timer.d/refresh-less-often.conf
[Timer]

# run less often (compared to distro-default of twice a day)
OnCalendar=monthly
